# services/backgroundchange/background.service
[Unit]
Description=Change the background

[Service]
Type=oneshot
ExecStart=%h/.config/systemd/user/backgroundchange
